# services/xorg/run
#!/data/data/com.termux/files/usr/bin/sh

exec 2>&1

exec /data/data/com.termux/files/usr/opt/quine/scripts/start-xorg.sh